# supervisord.conf
; Put this file in the path: /etc/supervisord.conf and run in a terminal: supervisord -c /etc/supervisord.conf

[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file
 
[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
 
;[program:AquaGame_8000]
;command = python /home/weblord/Desktop/AquaGame/manage.py 0.0.0.0:8000
;stderr_logfile = /var/log/supervisord/AquaGame_8000_stderr.log
;stdout_logfile = /var/log/supervisord/AquaGame_8000_stdout.log
;autorestart = true
 
;[program:RELab_8008]
;command = python /home/weblord/Desktop/RELab/manage.py 0.0.0.0:8008
;stderr_logfile = /var/log/supervisord/RELab_8008_stderr.log
;stdout_logfile = /var/log/supervisord/RELab_8008_stdout.log
;autorestart = true

[program:RoMIE_80]
command = python /home/weblord/Desktop/RoMIE/manage.py 0.0.0.0:80
stderr_logfile = /var/log/supervisord/RoMIE_80_stderr.log
stdout_logfile = /var/log/supervisord/RoMIE_80_stdout.log
autorestart = true

